# site-cookbooks/regulus/templates/default/rc.local.erb
/usr/local/rvm/rubies/ruby-<%= node[:regulus][:ruby_version] %>/bin/ruby <%= node[:regulus][:deploy_dir] %>/current/lib/backup.rb
