# site-cookbooks/alterf/templates/default/rc.local.erb
touch /var/lock/subsys/local
setterm -blank 0
/usr/local/rvm/rubies/ruby-<%= node[:alterf][:ruby_version] %>/bin/ruby <%= node[:alterf][:deploy_dir] %>/current/lib/backup.rb
